# mixins/configs/process/supervisord-port.conf
[supervisord]
nodaemon=true
logfile=%(ENV_LOG_PATH)s/supervisord-port.log
childlogdir=%(ENV_LOG_PATH)s

[program:android_port_forward]
command=/usr/local/bin/docker-android start port_forwarder
autostart=true
autorestart=true
stdout_logfile=%(ENV_LOG_PATH)s/android_port_forward.stdout.log
redirect_stderr=true
priority=2
